# packages/dofast/dofast-1.2.14-py3-none-any.whl/dofast/data/qflask.conf
[program:qflask]
directory=/tmp/
command=/usr/local/bin/qflask
user=root
outostart=true
autorestart=true
redirect_stderr=true
stopasgroup=true
killasgroup=true